# ansible/firewall/roles/infra-dns/files/etc/systemd/system/unbound-blocklist.service
[Unit]
Description=Unbound blocklist update service
After=network.target

[Service]
ExecStart=/etc/unbound/scripts/update-blocklists.sh
