# ansible/firewall/roles/infra-dns/files/etc/systemd/system/unbound-roothints.service
[Unit]
Description=Unbound root.hints update service
After=network.target

[Service]
ExecStart=/etc/unbound/scripts/update-roothints.sh
